# deploy/systemd/multica-daemon.service
[Unit]
Description=Multica local agent runtime daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=HOME=/root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# This VM is a dedicated AI experiment sandbox. Allow autonomous Claude runs
# while the daemon is intentionally managed as root.
Environment=IS_SANDBOX=1
# Use an interactive shell so the daemon inherits the local Claude/Codex
# runtime bootstrap from /root/.bashrc without copying API keys into this unit.
ExecStart=/bin/bash -ic 'exec /usr/local/bin/multica daemon start --foreground --device-name openclaw.oxygent.org.cn --runtime-name "OpenClaw Runtime" --max-concurrent-tasks 3 --no-auto-update'
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
